# ops/systemd/dhi.service
[Unit]
Description=Dhi AI Agent Security
After=network.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/local/bin/dhi --level alert --ebpf-ssl-only --port 9090
Restart=always
RestartPreventExitStatus=73
RestartSec=5
WorkingDirectory=/var/log/dhi
UMask=0027
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectKernelModules=true
LockPersonality=true
RestrictSUIDSGID=true
ReadWritePaths=/var/log/dhi
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_SYS_ADMIN CAP_SYS_PTRACE CAP_NET_ADMIN
CapabilityBoundingSet=CAP_BPF CAP_PERFMON CAP_SYS_ADMIN CAP_SYS_PTRACE CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
